# Stable-Diffusion-ComfyUI-Desktop/config/ollama.conf
[program:ollama]
user=ollama
command=/usr/local/bin/ollama serve
autostart=true
autorestart=true
stderr_logfile=/var/log/ollama.err
stdout_logfile=/var/log/ollama.log
environment=CODENATION_ENV=prod,HOME=/var/lib/ollama,USER=ollama,OLLAMA_MODELS=/var/lib/ollama/models