# .devcontainer/config/redis-server.conf
[program:redis-server]
command=/usr/bin/redis-server
autostart=false
autorestart=true
stderr_logfile=/var/log/redis-server.err
stdout_logfile=/var/log/redis-server.log
environment=CODENATION_ENV=prod